# launchd/com.denise.recovery.reminders.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.denise.recovery.reminders</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/bsoper/.local/bin/denise-reminders</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/bsoper/Library/Application Support/DeniseRecovery/reminder-runner</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>REMINDER_CONFIG_PATH</key>
    <string>/Users/bsoper/Library/Application Support/DeniseRecovery/reminder-runner/reminders.json</string>
    <key>REMINDER_STATE_PATH</key>
    <string>/Users/bsoper/Library/Application Support/DeniseRecovery/reminder-runner/reminder-state.json</string>
    <key>REMINDER_WINDOW_MINUTES</key>
    <string>90</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/bsoper/Library/Application Support/DeniseRecovery/reminder-runner/reminders-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bsoper/Library/Application Support/DeniseRecovery/reminder-runner/reminders-launchd.err</string>
</dict>
</plist>
